# iso_profile/airootfs/etc/systemd/system/spider-ttycolors.service
[Unit]
Description=Spider's LAB — TTY Color Palette
After=systemd-vconsole-setup.service
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/local/bin/spider-ttycolors
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
